# .config/systemd/user/uxplay.service
[Unit]
Description=uxplay daemon

[Service]
ExecStart=uxplay
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
